# packaging/alpine/sproxy.init
#!/sbin/openrc-run

name="sproxy"
description="sproxy high-performance proxy server"
command="/usr/bin/sproxy"
command_args="-D --pidfile /run/sproxy.pid -c /etc/sproxy/sproxy.conf"
pidfile="/run/sproxy.pid"

depend() {
	need net
	after firewall
}
